# upstart/kubelet.conf
description "Kubelet"

start on runlevel [2345]
stop on runlevel [!2345]

limit nofile 20000 20000

kill timeout 30 # wait 30s between SIGTERM and SIGKILL.

pre-start script
  KUBELET=/usr/bin/$UPSTART_JOB
  if [ -f /etc/default/$UPSTART_JOB ]; then
    . /etc/default/$UPSTART_JOB
  fi
  if [ -f $KUBELET ]; then
    exit 0
  fi
exit 22
end script

script
  KUBELET=/usr/bin/$UPSTART_JOB
  CONFIG=""
  if [ -f /etc/default/kube ]; then
    . /etc/default/kube
  fi

  if [ -f /etc/default/$UPSTART_JOB ]; then
    . /etc/default/$UPSTART_JOB
  fi
  exec /usr/bin/kubelet \
       $KUBE_LOGTOSTDERR \
       $KUBE_LOG_LEVEL \
       $KUBELET_API_SERVER \
       $KUBELET_ADDRESS \
       $KUBELET_PORT \
       $KUBELET_HOSTNAME \
       $KUBE_ALLOW_PRIV \
       $KUBELET_ARGS
end script
