# upstart/kube-controller-manager.conf
description "Kubernetes Controller Manager"

start on runlevel [2345]
stop on runlevel [!2345]

limit nofile 20000 20000

kill timeout 30 # wait 30s between SIGTERM and SIGKILL.

pre-start script
  CM=/usr/bin/$UPSTART_JOB
  if [ -f /etc/default/$UPSTART_JOB ]; then
    . /etc/default/$UPSTART_JOB
  fi
  if [ -f $CM ]; then
    exit 0
  fi
exit 22
end script

script
  CM=/usr/bin/$UPSTART_JOB
  CONFIG=""
  if [ -f /etc/default/kube ]; then
    . /etc/default/kube
  fi

  if [ -f /etc/default/$UPSTART_JOB ]; then
    . /etc/default/$UPSTART_JOB
  fi
  exec /usr/bin/kube-controller-manager \
       $KUBE_LOGTOSTDERR \
       $KUBE_LOG_LEVEL \
       $KUBE_MASTER \
       $KUBE_CONTROLLER_MANAGER_ARGS
end script
