# upstart/kube-apiserver.conf
description "Kubernetes API Server"

start on runlevel [2345]
stop on runlevel [!2345]

limit nofile 20000 20000

kill timeout 30 # wait 30s between SIGTERM and SIGKILL.

pre-start script
  APISERVER=/usr/bin/$UPSTART_JOB
  if [ -f /etc/default/$UPSTART_JOB ]; then
    . /etc/default/$UPSTART_JOB
  fi
  if [ -f $APISERVER ]; then
    exit 0
  fi
exit 22
end script

script
  APISERVER=/usr/bin/$UPSTART_JOB
  CONFIG=""
  if [ -f /etc/default/kube ]; then
    . /etc/default/kube
  fi

  if [ -f /etc/default/$UPSTART_JOB ]; then
    . /etc/default/$UPSTART_JOB
  fi
  exec /usr/bin/kube-apiserver \
       $KUBE_LOGTOSTDERR \
       $KUBE_LOG_LEVEL \
       $KUBE_ETCD_SERVERS \
       $KUBE_API_ADDRESS \
       $KUBE_API_PORT \
       $KUBELET_PORT \
       $KUBE_ALLOW_PRIV \
       $KUBE_SERVICE_ADDRESSES \
       $KUBE_ADMISSION_CONTROL \
       $KUBE_API_ARGS
end script
